# packaging/reoftpd-cleanup.timer
[Unit]
Description=Run reoftpd retention sweep daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
