# packaging/reoftpd-cleanup.service
[Unit]
Description=reoftpd retention sweep

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reoftpd cleanup --once --config /etc/reoftpd/reoftpd.toml
User=reoftpd
Group=reoftpd
ReadWritePaths=/srv/reolink
